# server/fortress-yfinance.service
[Unit]
Description=Fortress YFinance Bar Updater (GC/SI/ES/NQ every 5 min)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=fortress
WorkingDirectory=/opt/fortress/04_codebase
EnvironmentFile=/opt/fortress/.env
ExecStart=/opt/fortress/venv/bin/python -X utf8 tick_yfinance_updater.py --loop --interval 300
Restart=always
RestartSec=60
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
